# supervisord.conf
[supervisord]
nodaemon=true
loglevel=DEBUG

[program:caddy2]
command=caddy run --config /etc/caddy/Caddyfile --adapter caddyfile

[program:v2ray]
command=/usr/bin/v2ray -config /etc/v2ray/config.json